# systemd/alphacota-trader.service
[Unit]
Description=AlphaCota Polymarket Trading Loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/alphacota

# Required environment variables in /root/alphacota/.env:
#   POLYMARKET_MODE=live           (or paper for simulation)
#   POLYGON_RPC_URL=https://...    (Alchemy or public RPC)
#   POLYMARKET_API_KEY=...
#   POLYMARKET_API_SECRET=...
#   POLYMARKET_PRIVATE_KEY_ENC=... (AES-256-GCM encrypted key from vault)
#   VAULT_PASSWORD=...             (PBKDF2 vault password)
#   OPENROUTER_API_KEY=...         (for AI probability estimation)
EnvironmentFile=/root/alphacota/.env
Environment=POLYMARKET_MODE=live

# Run preflight checks before starting the main loop
ExecStartPre=/usr/bin/python3 -m core.polymarket_preflight
ExecStart=/usr/bin/python3 -m core.polymarket_loop --mode=live

Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alphacota-trader

# Safety limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
